# crescent/crontab
*/15 * * * * /usr/local/bin/python /opt/crescent/manage.py search_unknown
