# twitter/crontab
*/15 * * * * /usr/local/bin/python /opt/twitter/main.py search
